# root/etc/services.d/xteve/run
#!/usr/bin/with-contenv bash

# Run App
cd /app
exec \
s6-setuidgid abc xteve -port=34400 -config=/config
